# web_app/simple.conf
[program:myserver]
command= gunicorn -w 4 -b 0.0.0.0:80 --log-file /home/ubuntu/gunicornlog.log run:app

[supervisord]
logfile=/home/ubuntu/supervisord.log
loglevel=debug
user=root
